# conf/supervisord/webreg_celery_worker.ini
[program:webreg_celery_worker]
command=/home/webreg/env/bin/python /home/webreg/webreg/webreg/manage.py celery worker --pidfile='/home/webreg/pids/celery.pid' --logfile='/home/webreg/webreg/log/celery/celery.log' --loglevel=INFO --time-limit=300
user=webreg
stdout_logfile=/home/webreg/webreg/log/celery/worker.log
stderr_logfile=/home/webreg/webreg/log/celery/worker.err.log
directory=/home/webreg/webreg/webreg
environment=
    PYTHONPATH="/home/webreg/webreg/webreg",
    PATH=/home/webreg/env/bin:/usr/local/bin:/usr/bin:/usr/local/sbin:/usr/sbin:/usr/cachesys/bin:/home/webreg/.local/bin:/home/webreg/bin,
    LD_LIBRARY_PATH=/lib64:/home/webreg/cachesys/bin
killasgroup=true
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=600
