# conf/supervisord/webreg_gunicorn.ini
[program:webreg_gunicorn]
command=/home/webreg/env/bin/gunicorn config.wsgi:application --bind 127.0.0.1:43800 --workers=8 --pid=/home/webreg/pids/gunicorn.pid --log-file /home/webreg/webreg/log/gunicorn/gunicorn.log --timeout 300
user=webreg
directory=/home/webreg/webreg/webreg
umask=022
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
environment=
    PATH=/home/webreg/env/bin:/usr/local/bin:/usr/bin:/usr/local/sbin:/usr/sbin:/usr/cachesys/bin:/home/webreg/.local/bin:/home/webreg/bin,
    LD_LIBRARY_PATH=/lib64:/home/webreg/cachesys/bin
